# chef/teamcity/templates/default/teamcity-agent.service.erb
## THIS FILE IS MANAGED BY CHEF. ALL CHANGES WILL BE OVERWRITTEN!

[Unit]
Description=TeamCity Agent
Documentation=https://confluence.jetbrains.com/display/TCD10/
After=docker.service network.target syslog.service
Requires=docker.service network.target syslog.service

[Service]
ExecStart=<%= @system_dir %>/bin/agent.sh start
ExecStop=<%= @system_dir %>/bin/agent.sh stop
KillMode=process
RemainAfterExit=yes
Type=oneshot
User=<%= @user %>
SuccessExitStatus=0 143
WorkingDirectory=<%= @system_dir %>

[Install]
WantedBy=multi-user.target
